# supervisord/src/supervisord.conf
[supervisord]
nodaemon=true

[program:sleep_n_kill]
command=sleep 10 && killall ping && echo done

[program:ping]
command=ping web.de

